# cookbooks/bcpc/templates/default/upstart-kibana.conf.erb
################################################
#
#              Generated by Chef
#
################################################

description     "kibana log search interface"

start on virtual-filesystems
stop on runlevel [06]

# Respawn it if the process exits
respawn
respawn limit 5 30
expect fork

chdir /opt/kibana
setuid <%="#{node[:bcpc][:kibana][:user]}"%>
setgid <%="#{node[:bcpc][:kibana][:group]}"%>
console log

script
  exec ruby -r bundle/bundler/setup.rb kibana.rb < /dev/null >> /var/log/kibana/kibana.log &
  emit kibana-running
end script
